# project/config/development/supervisor/mock_crontab.conf
[program:mock_crontab]
command=watch -n 60 /usr/bin/php /var/www/spider_frame/public/cli.php spider:trigger
environment= ENV="development";
stopwaitsecs=5
restart=always
user=root
